# dashboard/eth-dashboard.service
[Unit]
Description=ETH Trading Dashboard (Flask API + Telegram bot)
Documentation=https://github.com/hedgeme/analytical-dashboards
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tecviva
Group=tecviva
WorkingDirectory=/home/tecviva/analytical-dashboards/dashboard
EnvironmentFile=/home/tecviva/analytical-dashboards/dashboard/.env
ExecStart=/home/tecviva/analytical-dashboards/dashboard/venv/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eth-dashboard

# Prevent the service from being killed on high memory
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
